# deploy/linux/regime-locus-master.service
# Regime-Locus-Matrix — systemd unit (template: use install-systemd.sh or replace @INSTALL_ROOT@).
#
# Prereqs on the VPS:
#   - Python 3.11+ venv at @INSTALL_ROOT@/.venv with: pip install -e ".[ibkr]" (and Massive deps as needed)
#   - Repo at @INSTALL_ROOT@, .env with IBKR_HOST/PORT/CLIENT_ID, Massive keys, etc.
#   - IBKR: TWS is GUI-only; use IB Gateway in headless/API mode on the VPS, or tunnel API to a machine
#     that runs Gateway. The master loop talks to whatever IBKR_HOST points at.
#   - This service runs run_master.py (no Tk GUI). For pnl_chart_dashboard.py use a local PC or Xvfb.
#
# Logs: journalctl -u regime-locus-master -f

[Unit]
Description=Regime-Locus-Matrix master loop (universe pipeline + paper equity + monitor)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@INSTALL_ROOT@
Environment=PYTHONUNBUFFERED=1
# systemd parses KEY=VALUE only; keep .env free of "export " prefixes
EnvironmentFile=-@INSTALL_ROOT@/.env
ExecStart=@INSTALL_ROOT@/.venv/bin/python scripts/run_master.py
Restart=on-failure
RestartSec=20
# Allow child processes (subprocess chains) time to exit
KillMode=mixed
TimeoutStopSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=regime-locus-master

[Install]
WantedBy=multi-user.target
