# deploy/linux/regime-locus-crew.service
# Regime-Locus-Matrix — Hermes AI crew (health + market context + commander)
#
# Prereqs on the VPS:
#   - Python venv for ExecStart (Hermes must be installed into THAT interpreter's site-packages.)
#     Many VPSes use /opt/rlm-venv/bin/python — then run from repo root:
#       /opt/rlm-venv/bin/pip install -e ".[hermes]"
#     If ExecStart is @INSTALL_ROOT@/.venv/bin/python (PEP 668: never use system pip):
#       cd @INSTALL_ROOT@ && .venv/bin/pip install -e ".[hermes]"
#     One-shot helper: scripts/setup_hermes_on_vps.sh (prefers /opt/rlm-venv when present).
#
#   LLM backend — choose ONE of the two options below and add to .env:
#
#   Option A: Groq cloud (free tier, fast — recommended for VPS without GPU)
#     GROQ_API_KEY=gsk_...            ← free key from console.groq.com
#     RLM_HERMES_MODEL=llama-3.3-70b-versatile   ← optional; this is the auto-default
#     (RLM_HERMES_BASE_URL and RLM_HERMES_API_KEY are set automatically when GROQ_API_KEY is present)
#
#   Option B: local Ollama
#     RLM_HERMES_BASE_URL=http://127.0.0.1:11434/v1
#     RLM_HERMES_API_KEY=ollama
#     RLM_HERMES_MODEL=llama3.2
#
#   Also add (Hermes channel bot; dedicated token):
#     RLM_HERMES_TELEGRAM_BOT_TOKEN=...
#     RLM_HERMES_TELEGRAM_CHAT_ID=...
#   Legacy fallback still supported:
#     TELEGRAM_BOT_TOKEN=...
#     TELEGRAM_NOTIFY_CHAT_ID=...
#
# Install:
#   cp deploy/linux/regime-locus-crew.service /etc/systemd/system/
#   sed -i 's|@INSTALL_ROOT@|/opt/Regime-Locus-Matrix|g' /etc/systemd/system/regime-locus-crew.service
#   systemctl daemon-reload
#   systemctl enable --now regime-locus-crew
#
# Logs: journalctl -u regime-locus-crew -f

[Unit]
Description=RLM Hermes crew — health facts, market context, commander (gate_state.json)
After=network-online.target regime-locus-master.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@INSTALL_ROOT@
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=@INSTALL_ROOT@/src:@INSTALL_ROOT@
EnvironmentFile=-@INSTALL_ROOT@/.env

# Intervals (seconds) — override here or in .env
Environment=CREW_HEALTH_INTERVAL=120
Environment=CREW_ANALYSIS_INTERVAL=300
Environment=CREW_BRIEFING_INTERVAL=600
# Scotty: auto-restart dead RLM systemd units (never restarts regime-locus-crew from inside this process)
Environment=SCOTTY_AUTO_RESTART=1
Environment=SCOTTY_RESTART_COOLDOWN_SEC=180

# After install, replace with /opt/rlm-venv/bin/python if that is your venv.
ExecStart=@INSTALL_ROOT@/.venv/bin/python scripts/run_crew.py

Restart=always
RestartSec=30
KillMode=mixed
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=regime-locus-crew

[Install]
WantedBy=multi-user.target
