# deploy/linux/regime-locus-master.service
# Regime-Locus-Matrix — systemd unit (template: use install-systemd.sh or replace @INSTALL_ROOT@).
#
# Prereqs on the VPS:
#   - Python 3.11+ venv at @INSTALL_ROOT@/.venv with: pip install -e ".[ibkr]" (and Massive deps as needed)
#   - Repo at @INSTALL_ROOT@, .env with IBKR_HOST/PORT/CLIENT_ID, Massive keys, etc.
#   - IBKR: TWS is GUI-only; use IB Gateway in headless/API mode on the VPS (see deploy/systemd/ib-gateway.service.example),
#     or tunnel API to a machine that runs Gateway. The master loop talks to whatever IBKR_HOST points at.
#   - Telegram: run EITHER this unit with run_master.py --telegram-bot OR a standalone rlm-telegram.service —
#     never both (duplicate bots / double notifications). See deploy/rlm-master-telegram.service.example.
#   - Enable exactly ONE master trading loop: this unit OR rlm-master-telegram.service OR
#     rlm-master-trader.service (mutual Conflicts in each unit file).
#   - This service runs run_master.py (no Tk GUI). For pnl_chart_dashboard.py use a local PC or Xvfb.
#
# Logs: journalctl -u regime-locus-master -f

[Unit]
Description=Regime-Locus-Matrix master loop (universe pipeline + paper equity + monitor)
Conflicts=rlm-master-telegram.service rlm-master-trader.service
After=network-online.target docker.service ib-gateway.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@INSTALL_ROOT@
Environment=PYTHONUNBUFFERED=1
# systemd parses KEY=VALUE only; keep .env free of "export " prefixes
EnvironmentFile=-@INSTALL_ROOT@/.env
ExecStart=@INSTALL_ROOT@/scripts/systemd_exec_python.sh scripts/run_master.py
# "always" = come back after reboot (via enable) and after any process exit; use "on-failure" for crash-only
Restart=always
RestartSec=20
# Allow child processes (subprocess chains) time to exit
KillMode=mixed
TimeoutStopSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=regime-locus-master

[Install]
WantedBy=multi-user.target
