# deploy/linux/rlm-nightly-opt.timer
[Unit]
Description=Timer for RLM nightly hyperparameter optimization (Mon-Fri 17:30 America/New_York)

[Timer]
OnCalendar=Mon..Fri *-*-* 17:30:00 America/New_York
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
