# deploy/linux/rlm-market-close.timer
[Unit]
Description=RLM Market Close Timer (NYSE) - fire at 16:30 America/New_York Mon-Fri

[Timer]
OnCalendar=Mon..Fri *-*-* 16:30:00 America/New_York
Persistent=false

[Install]
WantedBy=timers.target
