# deploy/linux/rlm-weekly-calibrate.timer
[Unit]
Description=Timer for RLM weekly regime model calibration (Sun 20:00 America/New_York)

[Timer]
OnCalendar=Sun *-*-* 20:00:00 America/New_York
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
