# deploy/linux/rlm-market-open.timer
[Unit]
Description=RLM Market Open Timer (NYSE) - fire at 09:00 America/New_York Mon-Fri

[Timer]
OnCalendar=Mon..Fri *-*-* 09:00:00 America/New_York
Persistent=false

[Install]
WantedBy=timers.target
