# deploy/systemd/rlm-walkforward-universe.timer
[Unit]
Description=Daily RLM universe walk-forward OOS batch (Mon-Fri 17:45 America/New_York)

[Timer]
OnCalendar=Mon..Fri *-*-* 17:45:00 America/New_York
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
