# deploy/rlm-finetune.timer
[Unit]
Description=Timer for RLM Kronos Weekend Fine-Tuning (Sat 02:00 America/New_York)

[Timer]
OnCalendar=Sat *-*-* 02:00:00 America/New_York
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
