# deploy/rlm-finetune.service
[Unit]
Description=RLM Kronos Weekend Fine-Tuning
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/regime-locus
EnvironmentFile=/opt/regime-locus/.env
# Enforce CPU limits to prevent freezing the VPS
Environment="RLM_MAX_CPU_THREADS=1"
Environment="OMP_NUM_THREADS=1"
ExecStart=/opt/regime-locus/venv/bin/python3 /opt/regime-locus/scripts/wrapper_finetune.py
StandardOutput=append:/var/log/rlm-finetune.log
StandardError=append:/var/log/rlm-finetune.error.log
